# integrations/systemd/geist-home.service
# geist-home.service — the resident home appliance daemon (see DEPLOY.md).
#
# Install (adapt the three /home/germar paths to your box):
#   sudo cp geist-home.service /etc/systemd/system/
#   sudo systemctl daemon-reload && sudo systemctl enable --now geist-home
#
# The socket lives in the Home Assistant config volume so the HA container
# (and its geist_conversation component) reaches it at /config/geist.sock.

[Unit]
Description=geist home appliance daemon (--serve)
After=local-fs.target

[Service]
User=germar
WorkingDirectory=/home/germar/geist-public
Environment=GEIST_AGENT_TRACE=0
ExecStart=/home/germar/geist-public/bin/pi5/release/tools/geist-home --serve /home/germar/ha-config/geist.sock
Restart=on-failure
RestartSec=5
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
